# 50-external-keyboard.rules
SUBSYSTEM=="usb", ATTR{idVendor}=="04d9", ATTR{idProduct}=="0167", RUN+="/home/james/keyboard.sh"
